# resources/service/nginx/run
#!/bin/sh -e

echo "Starting nginx"
exec nginx -g "daemon off;error_log stderr;"
